# meta-hidav/recipes-systemd/nfs-utils/files/rpc.mountd.service
[Unit]
Description=RPC mountd
Requires=rpcbind.target

[Service]
EnvironmentFile=/etc/default/nfs
ExecStart=/usr/sbin/rpc.mountd --foreground $OPTS_RPC_MOUNTD
