# meta-hidav/recipes-systemd/nfs-utils/files/var-lib-nfs-rpc_pipefs.mount
[Unit]
Description=rpc_pipefs

[Mount]
What=rpc_pipefs
Where=/var/lib/nfs/rpc_pipefs
Type=rpc_pipefs
