# meta-hipos/recipes-systemd/nfs-utils/files/var-lib-nfs-rpc_pipefs.mount
[Unit]
Description=rpc_pipefs
Requires=rpcbind.target
After=rpcbind.target

[Mount]
What=rpc_pipefs
Where=/var/lib/nfs/rpc_pipefs
Type=rpc_pipefs
