# meta-hidav/recipes-hidav/rootfs-overlay/rootfs-overlay/mount-rootfs-overlay.service
[Unit]
Description=AUFS overlay file system mount over root-fs
DefaultDependencies=false
Before=rootfs-overlay.target systemd-journald.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/mount-rootfs-overlay.sh
ExecStart=/bin/systemctl daemon-reload
ExecStartPost=/bin/systemctl --no-block start multi-user.target

