# meta-hidav/recipes-hidav/rootfs-overlay/rootfs-overlay/rootfs-overlay.target
[Unit]
Description=Default target
DefaultDependencies=false

