# meta-hidav/recipes-systemd/portmap/files/portmap.service
[Unit]
Description=port mapper service
Requires=network.target

[Service]
Type=forking
ExecStart=/sbin/portmap
Sockets=portmap.socket
StandardOutput=null

[Install]
WantedBy=multi-user.target

